# package/base-files-tiny/files/etc/init.d/done
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=95
boot() {
	[ -f /etc/rc.local ] && {
		sh /etc/rc.local
	}
}
